# scripts/com.utcomputational.gitlab-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.utcomputational.gitlab-export</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ben/Projects/UT_Computational_NE/Neutron_OS/tools/run_gitlab_export.sh</string>
    </array>
    
    <!-- Run every Sunday at 8:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Also run on load (catch up if missed) -->
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/ben/Projects/UT_Computational_NE/Neutron_OS/tools/exports/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/ben/Projects/UT_Computational_NE/Neutron_OS/tools/exports/launchd_stderr.log</string>
</dict>
</plist>
