# src/neutron_os/extensions/builtins/eve_agent/com.neutronos.signal-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neutronos.signal-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ben/Projects/UT_Computational_NE/.venv/bin/python</string>
        <string>-m</string>
        <string>tools.extensions.builtins.sense.cli</string>
        <string>serve</string>
        <string>--port</string>
        <string>8765</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ben/Projects/UT_Computational_NE/Neutron_OS</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/neutronos-signal.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/neutronos-signal.log</string>
</dict>
</plist>
